# lagoon/supervisord-worker-two.conf
[program:worker_two]
command=/app/start-n8n-queue-worker
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
stopwaitsecs=3600

